# config/etc/udev/rules.d/90-hiddev.rules
# SIPL 2014-05-19 Respond to hiddev events.

SUBSYSTEM=="usb", KERNEL=="hiddev*", RUN+="/etc/udev/scripts/hiddev.sh %n"
