# config/etc/init/evacs.conf
#
# This service starts the eVACS election server menu.

start on stopped rc RUNLEVEL=[2345]

stop on runlevel [016]

console owner

respawn
exec /opt/eVACS/bin/menu_start.sh
